# crontab
*/5 * * * * /usr/local/bin/python3 /app/manage.py crawler >> /var/log/cron.log 2>&1
